# provision/roles/openresty/files/nginx.conf
start on [2345]
stop  on [!2345]

exec /usr/sbin/nginx
